# templates/default/rootfs/chef-install.conf.erb
start on filesystem
exec /usr/sbin/install-chef.sh
