# templates/default/rootfs/init-net.conf.erb
start on startup

# we just fake the network coming up, lxc did this for us

emits net-device-up

script 
        <%- if bridge = @host[:bridges][    (@guest[:bridge] || 'vmbr0')    ] %>
        route add default gw <%= bridge[:ipv4][:address] %>
          <%- if bridge[:ipv6] %>
          ip -6 route add <%= bridge[:ipv6][:address] %> dev eth0
          ip -6 route add default via <%= bridge[:ipv6][:address] %> 
          <%- end %>
        <%- end %>
        /sbin/initctl emit --no-wait net-device-up IFACE=lo
        /sbin/initctl emit --no-wait net-device-up IFACE=eth0
end script
