# templates/default/rootfs/init-vm.conf.erb
start on startup

# we just fake the filesystem, lxc did this for us

emits filesystem

script 
        /sbin/initctl emit --no-wait filesystem
end script
